# etc/udev/rules.d/63-usb-optimization.rules
# ============================================================================
# USB Device Optimization - Better latency and stability
# ============================================================================
# Optimizes USB devices for responsiveness without sacrificing stability

# USB mice: disable autosuspend (can cause lag)
ACTION=="add", SUBSYSTEM=="usb", ATTRS{bInterfaceClass}=="03", ATTRS{bInterfaceSubClass}=="01", \
  ATTR{power/autosuspend}="-1"

# USB keyboards: keep powered
ACTION=="add", SUBSYSTEM=="usb", ATTRS{bInterfaceClass}=="03", ATTRS{bInterfaceSubClass}=="01", \
  ATTR{power/autosuspend}="-1"

# USB storage: prevent unexpected disconnects
ACTION=="add", SUBSYSTEM=="usb", ATTRS{bInterfaceClass}=="08", \
  ATTR{power/autosuspend}="10"

# Game controllers: real-time priority for input
ACTION=="add", SUBSYSTEM=="input", ATTRS{name}=="*Joypad*", \
  RUN+="/bin/bash -c 'echo 99 > /proc/$$(pidof -s Xvfb || pidof -s sway)/oom_score_adj 2>/dev/null || true'"

# Webcams: balanced power saving
ACTION=="add", SUBSYSTEM=="video4linux", SUBSYSTEMS=="usb", \
  ATTR{power/autosuspend_delay_ms}="5000"
