# etc/systemd/system/cpu-boost.service
[Unit]
Description=Intel CPU Boost for Responsiveness During I/O
Documentation=https://wiki.archlinux.org/title/CPU_frequency_scaling
After=multi-user.target

[Service]
Type=oneshot
User=root
ExecStart=/bin/bash -c 'echo 0 > /sys/devices/system/cpu/intel_pstate/no_turbo'
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cpu-boost

[Install]
WantedBy=multi-user.target
